# claude/lazycortex-review/skills/lazy-review.configure/templates/systemd.service
[Unit]
Description=lazycortex-review dispatcher (${REPO_ROOT})

[Service]
Type=oneshot
WorkingDirectory=${REPO_ROOT}
ExecStart=/usr/bin/env claude --agent lazy-review.dispatcher
